# .config/systemd/user/greenclip.service
[Unit]
Description=Greenclip agent
After=display-manager.service

[Service]
Type=simple
ExecStart=/home/g/.local/bin/greenclip daemon
Restart=always

[Install]
WantedBy=default.target
